# systemd/quotahalo-refresh.timer
[Unit]
Description=Refresh QuotaHalo usage cache every 30 seconds

[Timer]
OnActiveSec=1s
OnUnitActiveSec=30s
AccuracySec=1s
Unit=quotahalo-refresh.service
Persistent=false

[Install]
WantedBy=timers.target
